# scripts/english_bot.service
[Unit]
Description=English Learning Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/english_bot
EnvironmentFile=/opt/english_bot/.env
ExecStart=/opt/english_bot/venv/bin/python /opt/english_bot/bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
